# src/confs/hotspot-startup-script.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/hotspot-startup-script.sh

[Install]
WantedBy=default.target